# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:backend]
command=/bin/sh -lc 'if [ "${POSEIDON_MEMRAY_ENABLED:-0}" = "1" ]; then exec /usr/local/bin/memray run --force -o "${POSEIDON_MEMRAY_OUTPUT:-/app/backend/data/memray/memory_profile.bin}" -m uvicorn src.api.app:create_app --factory --host 127.0.0.1 --port 8000 --proxy-headers; else exec /usr/local/bin/uvicorn src.api.app:create_app --factory --host 127.0.0.1 --port 8000 --proxy-headers; fi'
directory=/app/backend
user=poseidon
autostart=true
autorestart=true
startsecs=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=0
stopsignal=QUIT
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20
